# roles/pihole/templates/pihole-config-watch.path.j2
[Unit]
Description=Watch PiHole config for changes
After=docker.service

[Path]
PathChanged={{ pihole_etc_dir }}/pihole.toml
Unit=pihole-backup-flag.service

[Install]
WantedBy=multi-user.target
